# deploy/ozon-parser.service
[Unit]
Description=Ozon Position Parser
After=network.target

[Service]
Type=oneshot
User=ozon
Group=ozon
WorkingDirectory=/opt/ozon-parser
Environment="PYTHONPATH=/opt/ozon-parser"
ExecStart=/opt/ozon-parser/.venv/bin/python app/main.py

# Logging
StandardOutput=append:/var/log/ozon-parser/parser.log
StandardError=append:/var/log/ozon-parser/parser.log

# Restart on failure
Restart=no

[Install]
WantedBy=multi-user.target
